# build/unison.ini
[program:unison]
command = unison -numericids -repeat watch -auto -batch /host_sync /app_sync -logfile /tmp/unison.log
redirect_stderr = true
autorestart = true
