# build/unison.ini
[program:unison]
command = unison default /host_sync /app_sync
#command = unison -ignore "Name {.git,.cache,*.swp}" -numericids -repeat watch -auto -batch /host_sync /app_sync -logfile /tmp/unison.log
redirect_stderr = true
autorestart = true
